# deploy/hostinger-vps/systemd/primegaming-scheduler.service
[Unit]
Description=PRIME Gaming Laravel Scheduler
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/gaming-platform/backend
ExecStart=/usr/bin/php artisan schedule:work
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target